# etc/services.d/fcgi/run
#!/usr/bin/execlineb -P

s6-setuidgid user
/usr/sbin/fcgiwrap -s unix:/tmp/fcgi.sock