# etc/services.d/webfused/run
#!/usr/bin/execlineb -P

s6-setuidgid user
/usr/local/bin/webfused -f /etc/webfused.conf
